# knowledge-graph/server/memory-mcp.service
[Unit]
Description=Knowledge Graph MCP Server
Documentation=https://github.com/mironmax/claudecode-plugins
After=network.target

[Service]
# kg-memory's start subcommand daemonizes via setsid and exits.
# RemainAfterExit keeps systemd's view of "active" tied to ExecStop running on shutdown.
Type=oneshot
RemainAfterExit=yes
ExecStart=%h/.local/bin/kg-memory start
ExecStop=%h/.local/bin/kg-memory stop
Restart=on-failure
RestartSec=5s

# Environment
Environment="KG_LOG_LEVEL=INFO"
Environment="KG_HTTP_HOST=127.0.0.1"
Environment="KG_HTTP_PORT=8765"

[Install]
WantedBy=default.target
